# var/lib/dpkg/info/inetutils-syslogd.conffiles
/etc/default/inetutils-syslogd
/etc/init.d/inetutils-syslogd
/etc/logrotate.d/inetutils-syslogd
/etc/syslog.conf
