# var/lib/dpkg/info/inetutils-syslogd.list
/.
/etc
/etc/default
/etc/default/inetutils-syslogd
/etc/init.d
/etc/init.d/inetutils-syslogd
/etc/logrotate.d
/etc/logrotate.d/inetutils-syslogd
/etc/syslog.conf
/etc/syslog.d
/usr
/usr/sbin
/usr/sbin/syslogd
/usr/share
/usr/share/doc
/usr/share/doc/inetutils-syslogd
/usr/share/doc/inetutils-syslogd/AUTHORS
/usr/share/doc/inetutils-syslogd/NEWS.gz
/usr/share/doc/inetutils-syslogd/THANKS
/usr/share/doc/inetutils-syslogd/changelog.Debian.gz
/usr/share/doc/inetutils-syslogd/changelog.gz
/usr/share/doc/inetutils-syslogd/copyright
/usr/share/man
/usr/share/man/man5/syslog.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/syslogd.8.gz
/var
/var/log
/var/log/news
