# var/lib/dpkg/info/tigervnc-standalone-server.list
/.
/etc
/etc/X11
/etc/X11/Xtigervnc-session
/etc/tigervnc
/etc/tigervnc/vncserver.users
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/tigervncserver@.service
/usr
/usr/bin
/usr/bin/Xtigervnc
/usr/bin/tigervncserver
/usr/libexec
/usr/libexec/tigervncsession-start
/usr/sbin
/usr/sbin/tigervncsession
/usr/share
/usr/share/doc
/usr/share/doc/tigervnc-standalone-server
/usr/share/doc/tigervnc-standalone-server/changelog.Debian.gz
/usr/share/doc/tigervnc-standalone-server/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/Xtigervnc.1.gz
/usr/share/man/man1/tigervncserver.1.gz
/usr/share/man/man8
/usr/share/man/man8/tigervncsession.8.gz
