# var/lib/dpkg/info/xrdp.list
/.
/etc
/etc/default
/etc/default/xrdp
/etc/init.d
/etc/init.d/xrdp
/etc/logrotate.d
/etc/logrotate.d/xrdp
/etc/pam.d
/etc/pam.d/xrdp-sesman
/etc/xrdp
/etc/xrdp/km-00000406.ini
/etc/xrdp/km-00000407.ini
/etc/xrdp/km-00000409.ini
/etc/xrdp/km-0000040a.ini
/etc/xrdp/km-0000040b.ini
/etc/xrdp/km-0000040c.ini
/etc/xrdp/km-00000410.ini
/etc/xrdp/km-00000411.ini
/etc/xrdp/km-00000412.ini
/etc/xrdp/km-00000414.ini
/etc/xrdp/km-00000415.ini
/etc/xrdp/km-00000416.ini
/etc/xrdp/km-00000419.ini
/etc/xrdp/km-0000041d.ini
/etc/xrdp/km-00000807.ini
/etc/xrdp/km-00000809.ini
/etc/xrdp/km-0000080a.ini
/etc/xrdp/km-0000080c.ini
/etc/xrdp/km-00000813.ini
/etc/xrdp/km-00000816.ini
/etc/xrdp/km-0000100c.ini
/etc/xrdp/km-00010409.ini
/etc/xrdp/km-19360409.ini
/etc/xrdp/pulse
/etc/xrdp/pulse/default.pa
/etc/xrdp/reconnectwm.sh
/etc/xrdp/sesman.ini
/etc/xrdp/startwm.sh
/etc/xrdp/xrdp.ini
/etc/xrdp/xrdp_keyboard.ini
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/xrdp-sesman.service
/lib/systemd/system/xrdp.service
/usr
/usr/bin
/usr/bin/xrdp-dis
/usr/bin/xrdp-genkeymap
/usr/bin/xrdp-keygen
/usr/bin/xrdp-sesadmin
/usr/bin/xrdp-sesrun
/usr/include
/usr/include/ms-erref.h
/usr/include/ms-fscc.h
/usr/include/ms-rdpbcgr.h
/usr/include/ms-rdpeclip.h
/usr/include/ms-rdpedisp.h
/usr/include/ms-rdpefs.h
/usr/include/ms-rdpegdi.h
/usr/include/ms-rdpele.h
/usr/include/ms-rdperp.h
/usr/include/ms-smb2.h
/usr/include/painter.h
/usr/include/rfxcodec_common.h
/usr/include/rfxcodec_decode.h
/usr/include/rfxcodec_encode.h
/usr/include/xrdp_client_info.h
/usr/include/xrdp_constants.h
/usr/include/xrdp_rail.h
/usr/include/xrdp_sockets.h
/usr/lib
/usr/lib/arm-linux-gnueabihf
/usr/lib/arm-linux-gnueabihf/libpainter.a
/usr/lib/arm-linux-gnueabihf/librfxencode.a
/usr/lib/arm-linux-gnueabihf/pkgconfig
/usr/lib/arm-linux-gnueabihf/pkgconfig/libpainter.pc
/usr/lib/arm-linux-gnueabihf/pkgconfig/rfxcodec.pc
/usr/lib/arm-linux-gnueabihf/pkgconfig/xrdp.pc
/usr/lib/arm-linux-gnueabihf/xrdp
/usr/lib/arm-linux-gnueabihf/xrdp/libcommon.a
/usr/lib/arm-linux-gnueabihf/xrdp/libcommon.so.0.0.0
/usr/lib/arm-linux-gnueabihf/xrdp/libmc.a
/usr/lib/arm-linux-gnueabihf/xrdp/libmc.so
/usr/lib/arm-linux-gnueabihf/xrdp/libscp.a
/usr/lib/arm-linux-gnueabihf/xrdp/libscp.so.0.0.0
/usr/lib/arm-linux-gnueabihf/xrdp/libvnc.a
/usr/lib/arm-linux-gnueabihf/xrdp/libvnc.so
/usr/lib/arm-linux-gnueabihf/xrdp/libxrdp.a
/usr/lib/arm-linux-gnueabihf/xrdp/libxrdp.so.0.0.0
/usr/lib/arm-linux-gnueabihf/xrdp/libxrdpapi.a
/usr/lib/arm-linux-gnueabihf/xrdp/libxrdpapi.so.0.0.0
/usr/lib/arm-linux-gnueabihf/xrdp/libxup.a
/usr/lib/arm-linux-gnueabihf/xrdp/libxup.so
/usr/sbin
/usr/sbin/xrdp
/usr/sbin/xrdp-chansrv
/usr/sbin/xrdp-sesman
/usr/share
/usr/share/bug
/usr/share/bug/xrdp
/usr/share/bug/xrdp/control
/usr/share/doc
/usr/share/doc/xrdp
/usr/share/doc/xrdp/NEWS.Debian.gz
/usr/share/doc/xrdp/README.Debian
/usr/share/doc/xrdp/changelog.Debian.gz
/usr/share/doc/xrdp/changelog.gz
/usr/share/doc/xrdp/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/xrdp
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/xrdp-dis.1.gz
/usr/share/man/man5/sesman.ini.5.gz
/usr/share/man/man5/xrdp.ini.5.gz
/usr/share/man/man8
/usr/share/man/man8/xrdp-chansrv.8.gz
/usr/share/man/man8/xrdp-genkeymap.8.gz
/usr/share/man/man8/xrdp-keygen.8.gz
/usr/share/man/man8/xrdp-sesadmin.8.gz
/usr/share/man/man8/xrdp-sesman.8.gz
/usr/share/man/man8/xrdp-sesrun.8.gz
/usr/share/man/man8/xrdp.8.gz
/usr/share/xrdp
/usr/share/xrdp/ad24b.bmp
/usr/share/xrdp/ad256.bmp
/usr/share/xrdp/cursor0.cur
/usr/share/xrdp/cursor1.cur
/usr/share/xrdp/sans-10.fv1
/usr/share/xrdp/socksetup
/usr/share/xrdp/xrdp24b.bmp
/usr/share/xrdp/xrdp256.bmp
/usr/share/xrdp/xrdp_logo.bmp
/usr/lib/arm-linux-gnueabihf/xrdp/libcommon.so
/usr/lib/arm-linux-gnueabihf/xrdp/libcommon.so.0
/usr/lib/arm-linux-gnueabihf/xrdp/libscp.so
/usr/lib/arm-linux-gnueabihf/xrdp/libscp.so.0
/usr/lib/arm-linux-gnueabihf/xrdp/libxrdp.so
/usr/lib/arm-linux-gnueabihf/xrdp/libxrdp.so.0
/usr/lib/arm-linux-gnueabihf/xrdp/libxrdpapi.so
/usr/lib/arm-linux-gnueabihf/xrdp/libxrdpapi.so.0
