# etc/profile.d/sudo.sh
[ -n "$PS1" -a "$(whoami)" = "android" ] || return 0
alias su="sudo su"
