# etc/update-motd.d/10-pihole
#!/bin/bash
[ -f /usr/local/bin/pihole ] && echo && bash -c "pihole version && echo && pihole status && echo"

